# centos6-cdh5-journalnode/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid
logfile=/var/log/supervisor/supervisord.log

[program:sshd]
command=service sshd start

[program:hadoop-hdfs-journalnode]
command=service hadoop-hdfs-journalnode start
