# centos6-cdh5-base/profile.d/oozie.sh
export OOZIE_URL=http://10.0.10.100:11000/oozie
